# dav/etc/systemd/system/mnt-dav.automount
[Unit]
Description=Mount perosnal WebDAV
After=network-online.target
Wants=network-online.target

[Automount]
Where=/mnt/dav
TimeoutIdleSec=300

[Install]
WantedBy=remote-fs.target
